# grafana/mitre_use_case.json
[{"mitre_tech_id":"T1037.004","mitre_det_id":"x-mitre-detection-strategy--be6a466c-40c6-4611-9b68-7cfcbcb35fb0","mitre_det_external_id":"DET0237","mitre_det_name":"Detection Strategy for Boot or Logon Initialization Scripts: RC Scripts","mitre_an_id":"x-mitre-analytic--69562961-14e6-42a7-9f8a-24ac00f6404e","mitre_an_external_id":"AN0660","mitre_an_name":"Analytic 0660","mitre_an_desc":"Detection of changes to \/etc\/rc.local.d\/local.sh or rc.local during post-boot script execution with abnormal commands or additions.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--685f917a-e95e-4ba0-ade1-c7d354dae6e0","name":"esxi:syslog","channel":"boot logs"},{"x_mitre_data_component_ref":"x-mitre-data-component--84572de3-9583-4c73-aabd-06ea88123dd8","name":"esxi:shell","channel":"admin command usage"}],"mitre_products":["esx","esx"]},{"mitre_tech_id":"T1037.004","mitre_det_id":"x-mitre-detection-strategy--be6a466c-40c6-4611-9b68-7cfcbcb35fb0","mitre_det_external_id":"DET0237","mitre_det_name":"Detection Strategy for Boot or Logon Initialization Scripts: RC Scripts","mitre_an_id":"x-mitre-analytic--98f8728d-ff74-47cb-b884-25071a21f77e","mitre_an_external_id":"AN0658","mitre_an_name":"Analytic 0658","mitre_an_desc":"Detection of modified or newly created \/etc\/rc.local or \/etc\/init.d scripts followed by suspicious execution during system startup.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--3d20385b-24ef-40e1-9f56-f39750379077","name":"auditd:SYSCALL","channel":"execve"},{"x_mitre_data_component_ref":"x-mitre-data-component--9f387817-df83-432a-b56b-a8fb7f71eedd","name":"linux:syslog","channel":"boot logs"}],"mitre_products":["linux","linux"]},{"mitre_tech_id":"T1037.004","mitre_det_id":"x-mitre-detection-strategy--be6a466c-40c6-4611-9b68-7cfcbcb35fb0","mitre_det_external_id":"DET0237","mitre_det_name":"Detection Strategy for Boot or Logon Initialization Scripts: RC Scripts","mitre_an_id":"x-mitre-analytic--b053dbd4-ad1e-45e1-a6b7-af2a5d931c82","mitre_an_external_id":"AN0661","mitre_an_name":"Analytic 0661","mitre_an_desc":"Detection of modified boot-time configuration scripts that persist malicious CLI commands across reboots.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--84572de3-9583-4c73-aabd-06ea88123dd8","name":"networkdevice:syslog","channel":"startup-config"},{"x_mitre_data_component_ref":"x-mitre-data-component--685f917a-e95e-4ba0-ade1-c7d354dae6e0","name":"networkdevice:syslog","channel":"system boot logs"}],"mitre_products":["network","network"]},{"mitre_tech_id":"T1037.004","mitre_det_id":"x-mitre-detection-strategy--be6a466c-40c6-4611-9b68-7cfcbcb35fb0","mitre_det_external_id":"DET0237","mitre_det_name":"Detection Strategy for Boot or Logon Initialization Scripts: RC Scripts","mitre_an_id":"x-mitre-analytic--e716b209-5b06-4bc4-843f-cbe4c51ddc0d","mitre_an_external_id":"AN0659","mitre_an_name":"Analytic 0659","mitre_an_desc":"Detection of edits or additions to \/etc\/rc.common, \/Library\/StartupItems, or \/System\/Library\/StartupItems and associated script execution during login or reboot.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--3d20385b-24ef-40e1-9f56-f39750379077","name":"macos:unifiedlog","channel":"process events"},{"x_mitre_data_component_ref":"x-mitre-data-component--2b3bfe19-d59a-460d-93bb-2f546adc2d2c","name":"fs:fsusage","channel":"file activity"}],"mitre_products":["macos","macos"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--0eb6cf59-4ba8-4cea-b64a-686ce7c69f70","mitre_an_external_id":"AN0266","mitre_an_name":"Analytic 0266","mitre_an_desc":"Use of native tools or scripting (e.g., `usermod`, `passwd`, `groupmod`) to escalate permissions or persist access on existing users, correlated with login or process events.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--d27b0089-2c39-4b6c-84ff-303e48657e77","name":"auditd:SYSCALL","channel":"usermod, groupmod, passwd"},{"x_mitre_data_component_ref":"x-mitre-data-component--84572de3-9583-4c73-aabd-06ea88123dd8","name":"auditd:PATH","channel":"\/etc\/passwd or \/etc\/group file write"}],"mitre_products":["linux","linux"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--5c69f3b9-8f73-455e-8eb1-5281cd6ce6d5","mitre_an_external_id":"AN0268","mitre_an_name":"Analytic 0268","mitre_an_desc":"Modifications to SSO\/SAML user attributes (e.g., `isAdmin`, `role`, MFA bypass, App assignments) often through CLI, API, or rogue IdP apps.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--d27b0089-2c39-4b6c-84ff-303e48657e77","name":"saas:okta","channel":"User Attribute Modified \/ Role Assignment Changed"}],"mitre_products":["saas"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--616ccbf4-08f2-4b54-8e41-a8e362e31827","mitre_an_external_id":"AN0267","mitre_an_name":"Analytic 0267","mitre_an_desc":"Modifications to user accounts via `dscl`, `pwpolicy`, or System Preferences CLI (`sysadminctl`) that alter user groups, enable root, or bypass MDM restrictions.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--d27b0089-2c39-4b6c-84ff-303e48657e77","name":"macos:unifiedlog","channel":"com.apple.accountsd, com.apple.opendirectoryd"}],"mitre_products":["macos"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--74565d24-df58-49b6-86e0-01a03d6dc2a7","mitre_an_external_id":"AN0269","mitre_an_name":"Analytic 0269","mitre_an_desc":"Addition of new users or changes to role permissions (e.g., ReadOnly -> Admin) via API or vSphere Client, particularly from non-jumpbox IPs.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--5b8b466b-2c81-4fe7-946f-d677a74ae3db","name":"esxi:vpxa","channel":"vim.SessionManager.login \/ vim.AccountManager.createUser"}],"mitre_products":["esx"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--842ba5ee-dcd0-42bd-9ef8-867a4ab1c703","mitre_an_external_id":"AN0265","mitre_an_name":"Analytic 0265","mitre_an_desc":"Account attribute changes (e.g., password set, group membership, servicePrincipalName, logon hours) correlated with unusual process lineage or timing, indicating privilege escalation or persistence via valid accounts.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--d27b0089-2c39-4b6c-84ff-303e48657e77","name":"WinEventLog:Security","channel":"EventCode=4738, 4728, 4670"},{"x_mitre_data_component_ref":"x-mitre-data-component--3d20385b-24ef-40e1-9f56-f39750379077","name":"WinEventLog:Sysmon","channel":"EventCode=1"}],"mitre_products":["windows","windows"]},{"mitre_tech_id":"T1098","mitre_det_id":"x-mitre-detection-strategy--d6c4cc3b-6875-4288-8193-bf4c864560ab","mitre_det_external_id":"DET0096","mitre_det_name":"Account Manipulation Behavior Chain Detection","mitre_an_id":"x-mitre-analytic--eb4a55f0-eff2-40f8-912e-43ba7e34603c","mitre_an_external_id":"AN0270","mitre_an_name":"Analytic 0270","mitre_an_desc":"Role escalation (e.g., Editor \u2192 Owner) in cloud collaboration tools (Google Workspace, O365) or file sharing apps to maintain elevated access.","mitre_log_sources":[{"x_mitre_data_component_ref":"x-mitre-data-component--d27b0089-2c39-4b6c-84ff-303e48657e77","name":"m365:unified","channel":"Admin Activity > Role Change or Sharing Change"}],"mitre_products":["m365"]}]